# packer/files/jumppad-connector.service
[Unit]
Description=Jumppad Connector

[Service]
WorkingDirectory=/root
ExecStart=/usr/local/bin/start-connector

[Install]
WantedBy=multi-user.target